# trunk/DEBIAN/postrm
#!/bin/sh
rm /etc/init.d/abyle
update-rc.d abyle remove
